# init/systemd/fproxy.service
[Unit]
Description=fproxy

[Service]
EnvironmentFile=-/etc/default/fproxy
ExecStart=/usr/local/bin/fproxy $FLAGS
User=fproxy
